# deb/rootfs/etc/services.d/ha-intergration/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Restarts the STDIN service
# Unregisters Service
# ==============================================================================

if [ "${1}" -gt 0 ]; then
  exec /run/s6/basedir/bin/halt
else
  bashio::log.info 'ha intergartion was stopped, restarting...'
fi
# vim: ft=sh
